# scripts/systemd/wayfi-hostapd.service
[Unit]
Description=WayFi Access Point (hostapd)
After=network.target
Before=wayfi-orchestrator.service

[Service]
Type=simple
ExecStart=/usr/sbin/hostapd /etc/hostapd/hostapd.conf
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
