# scripts/systemd/wayfi-llm.service
[Unit]
Description=WayFi LLM Server (llama.cpp)
After=network.target
Before=wayfi-orchestrator.service

[Service]
Type=simple
ExecStart=/usr/local/bin/llama-server \
    --model /opt/wayfi/models/llama-3.1-8b-instruct-q4_k_m.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    --ctx-size 2048 \
    --threads 4 \
    --batch-size 512
Restart=on-failure
RestartSec=10
MemoryMax=6G
MemoryHigh=5G

[Install]
WantedBy=multi-user.target
